# lib/systemd/system/openvpn@openvpn.service.d/50-unpriv.conf
## This file is part of Whonix.
## Copyright (C) 2012 - 2015 Patrick Schleizer <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Unit]
ConditionPathExists=!/var/run/qubes-service/whonix-template

[Service]

User=tunnel
Group=tunnel

ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
RestrictRealtime=true
ProtectKernelModules=true
CapabilityBoundingSet=
CapabilityBoundingSet=CAP_NET_ADMIN
ReadWritePaths=/etc/resolv.conf /run/openvpn-client /run/resolvconf
ProtectKernelTunables=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
SystemCallFilter=~@clock @cpu-emulation @debug @keyring @module @mount @obsolete @raw-io
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes
AmbientCapabilities=CAP_NET_ADMIN
