# lib/systemd/system/openvpn@openvpn.service.d/50-unpriv.conf
## This file is part of Whonix.
## Copyright (C) 2012 - 2015 Patrick Schleizer <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Unit]
ConditionPathExists=!/var/run/qubes-service/whonix-template

[Service]

User=tunnel

ExecStartPre=/usr/bin/sudo /usr/sbin/openvpn --rmtun --dev tun0
ExecStartPre=/usr/bin/sudo /usr/sbin/openvpn --mktun --dev tun0 --dev-type tun --user tunnel --group tunnel

ExecStopPost=/usr/bin/sudo /usr/sbin/openvpn --rmtun --dev tun0
